# deploy/systemd/regulon-agent.service
[Unit]
Description=Sannidh Live Regulatory Ingestion Agent
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/sannidh
Environment=NODE_ENV=production
Environment=ALERT_AGENT_PORT=8787
EnvironmentFile=/etc/sannidh/sannidh-agent.env
ExecStart=/usr/bin/node /opt/sannidh/index.js
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
StandardOutput=append:/var/log/sannidh-agent/out.log
StandardError=append:/var/log/sannidh-agent/err.log

[Install]
WantedBy=multi-user.target
